# scripts/templates/tasks.cron.jinja2
# Backup
00 3 * * *  {{ USER }}  bash {{ PWD }}/mc backup > {{ PWD }}/task-backup.log 2>&1

