# scripts/templates/tasks.cron.jinja2
# Backup
00 3 * * *  {{ USER }}  bash {{ PWD }}/mc backup

# Overviewer - generate map (incremental)
30 3 * * *  {{ USER }}  bash {{ PWD }}/mc map
